# IoT_Gateway/growhouse/gateway-db410c/apps_proc/sources/meta-growhouse/recipes-services/ble-advertise/files/ble-advertise.service
[Unit]
Description=Ble advertize service
After=bluetooth.target
Requires=bluetooth.target

[Service]
ExecStart=/etc/ble-advertise/ble-advertise.sh

[Install]
WantedBy=multi-user.target
