# IoT_Gateway/growhouse/gateway-db410c/apps_proc/sources/meta-growhouse/recipes-apps/zigbee-coordinator/files/zigbee-coordinator.service
[Unit]
Description=Zigbee coordinator application start service
After=mosquitto.service
Requires=mosquitto.service

[Service]
ExecStart=/etc/zigbee-coordinator/run_zigbee-coordinator.sh
RestartSec=5sec
Restart=always

[Install]
WantedBy=multi-user.target
